# files/etc/uci-defaults/02-disable-dhcp.sh
uci set dhcp.lan.ignore='1'
uci commit dhcp
