# files/etc/uci-defaults/03-mosdns.sh
uci set mosdns.config.enable='1'
uci set mosdns.config.listen='127.0.0.1:5335'
uci commit mosdns
