# files/etc/uci-defaults/04-passwall.sh
uci set passwall.global_forwarding.tproxy='1'
uci set passwall.global_forwarding.dns_mode='disabled'
uci commit passwall
